# satpy/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   

[supervisord]
logfile=/var/log/supervisord.log 
logfile_maxbytes=50MB        
logfile_backups=10          
loglevel=info              
pidfile=/tmp/supervisord.pid 
nodaemon=true               
childlogdir=/data/logs            
umask=022

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
environment=PYTHONUNBUFFERED=1

[program:nameserver]
command=/usr/local/bin/nameserver

[program:trollstalker]
command=/usr/local/bin/trollstalker.py -c /app/trollConfig/trollstalker_config.ini -C omps

[program:processOmps]
command=/app/bin/processOmps.py